# com.djtchill.streamflacr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.djtchill.streamflacr</string>
    <key>ProgramArguments</key>
    <array>
        <string>streamflacr</string>
        <string>--daemon</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/djtchill/Documents/Projects/StreamFLACr</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/djtchill/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/djtchill/Library/Logs/streamflacr.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/djtchill/Library/Logs/streamflacr.err</string>
</dict>
</plist>
